# cron/scripts/tests/mocks/createTransaction/crontab
#!/bin/sh

DIR="$(cd "$(dirname "$0")" && pwd)"
MOCKED_CRONTAB_ENTRIES="$DIR/mocked_crontab_entries.txt"

# Handle 'crontab -l' to list current entries
if [ "$1" = "-l" ]; then
    cat "$MOCKED_CRONTAB_ENTRIES"
    exit 0
fi

# Handle 'crontab <file>' or piped input to update entries
cat >"$MOCKED_CRONTAB_ENTRIES"
