# cron/scripts/tests/mocks/getPayrollsByEmployee/crontab
#!/bin/sh

DIR="$(cd "$(dirname "$0")" && pwd)"
MOCKED_CRONTAB_ENTRIES="$DIR/mocked_crontab_entries.txt"

case "$1" in
-l)
    cat "$MOCKED_CRONTAB_ENTRIES"
    ;;
*)
    if [ -f "$1" ]; then
        # If an argument is provided and it's a file, read from that file
        cat "$1" >"$MOCKED_CRONTAB_ENTRIES"
    else
        # Otherwise, capture the input being piped into the script
        cat >"$MOCKED_CRONTAB_ENTRIES"
    fi
    ;;
esac
